# linux/contrib/cfg2html.cron
#!/bin/sh
#
# /etc/cron.weekly/cfg2html

[ -x /usr/bin/cfg2html ] || exit 0

DIR=/var/log/cfg2html

if ! [ -d $DIR ]
then
    mkdir -p $DIR
    chown root:root $DIR
    chmod 0755 $DIR
fi

/usr/bin/cfg2html -o $DIR >/dev/null 2>&1

#



